# ansible/roles/pico-web/templates/ctf-daemon.service.j2
[Unit]
Description= ctf daemon runner
After=network.target vboxadd-service.service

[Service]
Type=simple
Environment="APP_SETTINGS_FILE={{ web_config_dir }}/deploy_settings.py"
ExecStartPre=/bin/sleep 15
ExecStart={{ virtualenv_dir }}/bin/daemon_manager -d {{ daemon_src_dir }} cache_stats share_instances
RestartSec="5min"
Restart=always

[Install]
WantedBy=multi-user.target
